# smartdns_config/alpine 用启动文件/smartdns
#!/sbin/openrc-run

name=$RC_SVCNAME
command="/opt/smartdns/smartdns"
command_args="-c /opt/smartdns/smartdns.conf"
pidfile="/run/$RC_SVCNAME.pid"
command_background="yes"

depend() {
        need net
}
